# docs/deploy/heritage.service
[Unit]
Description=Heritage Sites Finder API
After=network.target

[Service]
User=deploy
Group=deploy
WorkingDirectory=/srv/heritage
Environment=PATH=/srv/heritage/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/srv/heritage/.venv/bin/gunicorn --chdir /srv/heritage --bind 127.0.0.1:5000 "src.api:create_app()"
Restart=on-failure

[Install]
WantedBy=multi-user.target

# Notes:
# - Adjust User/Group and WorkingDirectory to match your deploy user and path.
# - This unit uses gunicorn to serve the Flask app; ensure gunicorn is in requirements.txt or installed in the venv.